# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
chmod +x /etc/auto-provision/stage2.sh
sh /etc/auto-provision/stage2.sh

exit 0
